# admin/com.futhark.build-runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.futhark.build-runner</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/martin/futhark/admin/build-runner.sh</string>
	</array>

	<key>WorkingDirectory</key>
	<string>/Users/martin/futhark</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>BUILD_REMOTE_HOST</key>
		<string>ubuntu@kenaz.app</string>
		<key>BUILD_REMOTE_PATH</key>
		<string>/var/www/kenaz.app/releases</string>
		<key>BUILD_REMOTE_HTML</key>
		<string>/var/www/kenaz.app</string>
	</dict>

	<key>StartInterval</key>
	<integer>600</integer>

	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/Users/martin/.futhark/build-runner.log</string>

	<key>StandardErrorPath</key>
	<string>/Users/martin/.futhark/build-runner.log</string>
</dict>
</plist>
